# backend/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/app/logs
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:uwsgi]
command=uwsgi --ini /app/docker/uwsgi.ini
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/uwsgi_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
priority=800

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/nginx_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
priority=900

[program:migrate]
command=python manage.py migrate --noinput
directory=/app
user=appuser
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/app/logs/migrate.log
stderr_logfile=/app/logs/migrate_error.log
priority=100

[program:collectstatic]
command=python manage.py collectstatic --noinput --clear
directory=/app
user=appuser
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/app/logs/collectstatic.log
stderr_logfile=/app/logs/collectstatic_error.log
priority=200